# image/stage2/04-set-timezone/files/tzupdate.service
[Unit]
Description=Timezone Update Service
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/tzupdate

[Install]
WantedBy=multi-user.target
